# slurmdbd/supervisor-conf/slurmdbd.conf
[program:runtime_agent]
command=/opt/prefix-output.sh /opt/runtime-agent.sh
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=100

[program:slurmdbd]
command=/opt/prefix-output.sh /opt/slurm/sbin/slurmdbd -D %(ENV_SLURMDBD_ARGS)s
autostart=true
autorestart=true
user=slurm
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
priority=150
